# cloudcontroller/files/salt-cloud.conf
[program:salt-cloud-dnsmasq]
command=dnsmasq -k -C /srv/scripts/dnsmasq.conf
redirect_stderr=true
